# deployment/systemd/colosseum-curator.service
[Unit]
Description=Colosseum CuratorAgent - The Keeper of Records
Documentation=https://github.com/davdunc/colosseum
After=colosseum-postgres.service
Requires=colosseum-postgres.service
Wants=network-online.target

[Service]
Type=simple
User=colosseum
Group=colosseum
WorkingDirectory=/opt/colosseum

# Environment
Environment="PATH=/opt/colosseum/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/colosseum/curator.env

# Database password from podman secret
# Database password from podman secret
ExecStartPre=/bin/bash -c 'mkdir -p -m 700 /run/user/%U/colosseum && podman secret inspect colosseum-db-password --showsecret > /run/user/%U/colosseum/db-password && chmod 600 /run/user/%U/colosseum/db-password'
Environment="DB_PASSWORD_FILE=/run/user/%U/colosseum/db-password"

# Start curator with default watchlist
ExecStart=/opt/colosseum/venv/bin/python -m colosseum.cli.curator start \
    --interval 60 \
    --tickers AAPL GOOGL MSFT AMZN TSLA NVDA META

# Cleanup
ExecStopPost=/bin/rm -f /run/user/%U/colosseum/db-password

# Restart policy
Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=colosseum-curator

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/colosseum/logs

# Resource limits
MemoryMax=1G
MemoryHigh=800M
CPUQuota=50%
TasksMax=50

[Install]
WantedBy=multi-user.target
